# murmel.service
[Unit]
Description=murmel – push-to-talk dictation daemon (Parakeet ASR)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/murmel start
Restart=on-failure
RestartSec=5

KillSignal=SIGTERM

Environment=WAYLAND_DISPLAY=wayland-1
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=graphical-session.target
